# assets/marxist-quote-fetch.service
[Unit]
Description=Fetch daily quote using getquotes
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/marxist_quote --fetch
